# ups/wol.service
[Unit]
Description=Python WOL service

[Service]
Type=simple
ExecStart=/home/darkacid/src/useful-stuff/ups/wol.py

[Install]
WantedBy=default.target
